# root/etc/s6-overlay/s6-rc.d/svc-homeassistant/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

setcap 'cap_net_bind_service=+ep' /usr/local/bin/python3.12

if [[ -z "${DISABLE_JEMALLOC+x}" ]]; then
  export LD_PRELOAD="/usr/local/lib/libjemalloc.so.2"
  export MALLOC_CONF="background_thread:true,metadata_thp:auto,dirty_decay_ms:20000,muzzy_decay_ms:20000"
fi

exec \
    s6-notifyoncheck -d -n 60 -w 5000 -c "nc -z localhost 8123" \
    s6-setuidgid abc python3 -m homeassistant -c /config
