# root/etc/s6-overlay/s6-rc.d/svc-homeassistant/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

setcap 'cap_net_bind_service=+ep cap_net_raw,cap_net_admin+eip' /usr/bin/python3.11

exec \
    s6-notifyoncheck -d -n 60 -w 5000 -c "nc -z localhost 8123" \
    s6-setuidgid abc /config/lsiopy/bin/python3 /lsiopy/bin/hass -c /config
