# root/etc/s6-overlay/s6-rc.d/init-config-homeassistant/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# Create secondary venv in /config for pip installs
if [[ ! -d /config/lsiopy/bin ]]; then
  mkdir -p /config/lsiopy
  python3 -m venv /config/lsiopy
fi
if ! grep -q '/config/lsiopy/bin' /var/run/s6/container_environment/PATH; then
  sed -i '1s|^|/config/lsiopy/bin:|' /var/run/s6/container_environment/PATH
fi

# set permissions
lsiown -R abc:abc \
    /config
