# root/etc/s6-overlay/s6-rc.d/init-config-homeassistant/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# Create secondary venv in /config for pip installs
if [[ ! -d /config/lsiopy/bin ]]; then
  s6-setuidgid abc python3 -m venv /config/lsiopy
fi
s6-setuidgid abc /config/lsiopy/bin/pip install \
    --no-cache-dir --upgrade \
    "pip>=21.0,<22.1" \
    setuptools \
    wheel

# set permissions
lsiown -R abc:abc \
    /config
