# root/etc/s6-overlay/s6-rc.d/init-config-homeassistant/run
#!/usr/bin/with-contenv bash

# set permissions
chown -R abc:abc \
    /config
