# root/etc/s6-overlay/s6-rc.d/init-config-homeassistant/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# migration
rm -rf /config/lsiopy

# set permissions
echo "Setting permissions"
lsiown -R abc:abc \
    /config
